# deployment/test-worker/runit/log/run
#!/bin/bash
LOG_PATH=/var/log/service/dpxdt_worker
mkdir -p $LOG_PATH
cd $LOG_PATH
/sbin/svlogd -tt $LOG_PATH
