# deployment/test-worker/runit/log/run
#!/bin/bash
export LOG_PATH=/var/log/service/dpxdt_worker
mkdir -p $LOG_PATH
exec /usr/bin/svlogd -ttt $LOG_PATH
